# systemd/xfce-theme-tray.service
[Unit]
Description=XFCE Theme Switch — System Tray
After=dbus.target

[Service]
Type=simple
ExecStartPre=/bin/sleep 3
ExecStart=/usr/bin/xfce-theme-tray
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
