# systemd/xfce-theme-tray.service
[Unit]
Description=XFCE Theme Switch — System Tray
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/xfce-theme-tray
Restart=on-failure
RestartSec=5

[Install]
WantedBy=graphical-session.target
