# systemd/xfce-theme-switch.timer
[Unit]
Description=Check XFCE theme every 10 minutes

[Timer]
OnBootSec=30s
OnUnitActiveSec=10min

[Install]
WantedBy=timers.target
