# systemd/xfce-theme-switch.service
[Unit]
Description=XFCE theme switch based on sun position

[Service]
Type=oneshot
ExecStart=/usr/bin/xfce-theme-switch
